# deploy/govfund-tracker.timer
[Unit]
Description=Run GovFund tracker ingestion daily

[Timer]
OnCalendar=*-*-* 06:00:00
Persistent=true
Unit=govfund-tracker.service

[Install]
WantedBy=timers.target
